# docker/worker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor/supervisor.sock   ; the path to the socket file
chown=www-data:www-data
chmod=0770                 ; socket file mode (default 0700)

[supervisord]
nodaemon=true
user=www-data
logfile=/var/log/supervisor/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB       ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10          ; # of main logfile backups; 0 means none, default 10
loglevel=info              ; log level; default info; others: debug,warn,trace
pidfile=/var/run/supervisord.pid ; supervisord pidfile; default supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:messenger-consume]
command=php /var/www/html/bin/console messenger:consume async --time-limit=3600 --limit=10 --memory-limit=128M
directory=/var/www/html
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
startsecs=0
autostart=true
autorestart=true
startretries=10
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=MESSENGER_CONSUMER_NAME=%(program_name)s_%(process_num)02d,APP_ENV=dev
